# profile.d/install.sh
sudo apt install xcape

sudo cp xcape.sh /etc/profile.d/xcape.sh
