# profile.d/xcape.sh
#!/bin/bash

xcape -t 100 -e 'Caps_Lock=Escape'
